# .config/systemd/user/compton.service
[Unit]
Description=Compton compositor

[Service]
ExecStart=/usr/bin/compton
Restart=always

[Install]
WantedBy=default.target
