# .config/systemd/user/polybar.service
[Unit]
Description=Polybar

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/polybar top
Restart=always

[Install]
WantedBy=default.target
