# .config/systemd/user/polybar.service
[Unit]
Description=Polybar

[Service]
ExecStart=/home/ole/.nix-profile/bin/polybar top
Restart=always

[Install]
WantedBy=default.target
