# .config/systemd/user/sxhkd.service
[Unit]
Description=Simple X Hotkey Daemon

[Service]
ExecStart=/nix/var/nix/profiles/per-user/ole/profile/bin/sxhkd -c /nix/var/nix/profiles/per-user/ole/profile/etc/sxhkd/sxhkdrc

[Install]
WantedBy=default.target
