# .config/systemd/user/sxhkd.service
[Unit]
Description=Simple X Hotkey Daemon
Requires=xorg.socket
After=xorg.socket

[Service]
ExecStart=/nix/var/nix/profiles/per-user/ole/profile/bin/sxhkd -c /nix/var/nix/profiles/per-user/ole/profile/etc/sxhkd/sxhkdrc
Restart=always

[Install]
WantedBy=default.target
